# resources/install-crictl.service
[Unit]
Description=Installs crictl
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/bin/install-crictl
[Install]
WantedBy=multi-user.target
