# resources/prometheus-ro-rootfs.service
[Unit]
Description=Runs ro-rootfs

[Service]
Type=oneshot
ExecStart=/opt/bin/prometheus-ro-rootfs
